# docker/config/etc/profile.d/conda.sh
#!/bin/bash

PATH=/opt/conda/bin:${PATH}
export PATH
source activate base
