# docker/config/etc/profile.d/cuda.sh
#!/bin/sh

PATH=/usr/local/cuda/bin:${PATH}
export PATH
