# docker/config/etc/profile.d/spack.sh
#!/bin/bash

if [ -d /opt/spack/bin ]; then
    PATH=/opt/spack/bin:${PATH}
    export PATH
fi
